# root/.config/systemd/user/goimapnotify.service
[Unit]
Description=GoIMAPNotify - IMAP IDLE mail watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PATH=%h/Code/go/bin:%h/go/bin:%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env goimapnotify
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
